# vpsproxy.ini
[program:proxy_service]
environment = LD_LIBRARY_PATH="/usr/local/lib:$LD_LIBRARY_PATH", PYTHONPATH="$PYTHONPATH:/usr/local/lib/python3.7/site-packages/vpsproxy"
command = /usr/local/bin/python3.7 /usr/local/bin/run_proxy_service
autostart = true
startsecs = 5
autorestart = true
startretries = 3
redirect_stderr = true
stdout_logfile_maxbytes = 20MB
stdout_logfile_backups = 20
stdout_logfile = /logs/vpsproxy/proxy_service.log

[program:broker_service]
environment = LD_LIBRARY_PATH="/usr/local/lib:$LD_LIBRARY_PATH", PYTHONPATH="$PYTHONPATH:/usr/local/lib/python3.7/site-packages/vpsproxy"
command = /usr/local/bin/python3.7 /usr/local/bin/run_broker_service
autostart = true
startsecs = 5
autorestart = true
startretries = 3
redirect_stderr = true
stdout_logfile_maxbytes = 20MB
stdout_logfile_backups = 20
stdout_logfile = /logs/vpsproxy/broker_service.log

[group:vpsproxy]
programs=proxy_service,broker_service
